# rootfs/etc/s6-overlay/s6-rc.d/chrony-config/up
/etc/s6-overlay/scripts/chrony-config.sh
